# prod/deploy/systemd/cff-voice-worker.service
[Unit]
Description=Content Fabric Voice Worker (voice change processing)
After=network.target redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/cff/prod
EnvironmentFile=/root/cff/.env
EnvironmentFile=/root/cff/prod/.env/.env.db
ExecStart=/root/cff/prod/venv/bin/python -m workers.voice_worker
Restart=always
RestartSec=5
StandardOutput=append:/var/log/cff-voice-worker.log
StandardError=append:/var/log/cff-voice-worker.log

[Install]
WantedBy=multi-user.target
