# prod/deploy/systemd/cff-scheduler.service
[Unit]
Description=Content Fabric Scheduler (polls DB for pending tasks)
After=network.target mysql.service redis.service
Wants=redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/cff/prod
EnvironmentFile=/root/cff/.env
EnvironmentFile=/root/cff/prod/.env/.env.db
EnvironmentFile=/root/cff/prod/.env/.env.api
ExecStart=/root/cff/prod/venv/bin/python -m scheduler.run
Restart=always
RestartSec=10
StandardOutput=append:/var/log/cff-scheduler.log
StandardError=append:/var/log/cff-scheduler.log

[Install]
WantedBy=multi-user.target
