# prod/deploy/systemd/cff-dle-processor-worker.service
[Unit]
Description=Content Fabric DLE Processor Worker (download + voice + assemble video)
After=network.target redis.service mysql.service
Wants=redis.service
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
User=root
WorkingDirectory=/opt/content-fabric/prod
EnvironmentFile=/opt/content-fabric/.env
EnvironmentFile=-/opt/content-fabric/prod/.env
ExecStart=/opt/content-fabric/prod/venv/bin/python -m workers.dle_processor_worker
Restart=always
RestartSec=5
StandardOutput=append:/var/log/cff-dle-processor-worker.log
StandardError=append:/var/log/cff-dle-processor-worker.log
SyslogIdentifier=cff-dle-processor
MemoryMax=2G
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
