# prod/deploy/systemd/cff-dle-ingestion-worker.service
[Unit]
Description=CFF DLE Ingestion Worker
After=network.target redis.service mysql.service
Wants=redis.service mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/content-fabric/prod
EnvironmentFile=/opt/content-fabric/.env
EnvironmentFile=-/opt/content-fabric/prod/.env
ExecStart=/opt/content-fabric/prod/venv/bin/python -m workers.dle_ingestion_worker
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cff-dle-ingestion

[Install]
WantedBy=multi-user.target
