# prod/deploy/systemd/cff-api.service
[Unit]
Description=Content Fabric API Gateway (FastAPI/uvicorn)
After=network.target mysql.service redis.service
Wants=redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/content-fabric/prod
EnvironmentFile=/opt/content-fabric/.env
EnvironmentFile=/opt/content-fabric/prod/.env/.env.db
EnvironmentFile=/opt/content-fabric/prod/.env/.env.api
ExecStart=/opt/content-fabric/prod/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
StandardOutput=append:/var/log/cff-api.log
StandardError=append:/var/log/cff-api.log

[Install]
WantedBy=multi-user.target
