# prod/deploy/systemd/cff-stats-worker.service
[Unit]
Description=CFF Stats Worker
After=network.target redis.service mysql.service
Wants=redis.service mysql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/content-fabric
Environment="PATH=/opt/content-fabric/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
EnvironmentFile=/opt/content-fabric/prod/.env/.env.db
EnvironmentFile=/opt/content-fabric/prod/.env/.env.api
ExecStart=/opt/content-fabric/venv/bin/python -m workers.stats_worker
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cff-stats

[Install]
WantedBy=multi-user.target
